# private_dot_config/private_autostart/private_blueman-tray.desktop
[Desktop Entry]
Name=Blueman Tray
Icon=blueman
Exec=blueman-tray
Terminal=false
Type=Application
StartupNotify=false
Categories=GTK;GNOME;Settings;HardwareSettings;
